# webapp/MASS-flask.service
[Unit]
Description=Wastewater management web app
After=network.target

[Service]
